# config/update-mellophone.timer
[Unit]
Description=Triggers mellophone-restart.service every day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
